# implementation/SpringBoot/src/main/resources/service/springboot.service
[Unit]
Description=springboot
After=syslog.target

[Service]
User=vagrant
ExecStart=/home/vagrant/springboot.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target