# implementation/SpringBoot/src/main/resources/service/springboot.1.service
[Unit]
Description=springboot
After=syslog.target

[Service]
User=afonsobarrenha
ExecStart=/home/afonsobarrenha/Documents/Projects/pocs/implementation/SpringBoot/target/springboot-2.0.0-SNAPSHOT.jar
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target